# 06_Linux_Process/Exercise_1/monitor.service
[Unit]
Description=Monitor Service
After=network.target

[Service]
ExecStart=/usr/local/bin/monitor
Restart=on-failure
RestartSec=5
StartLimitBurst=3
StartLimitIntervalSec=30

[Install]
WantedBy=multi-user.target
